# extra/supervisor.ini
[program:xstorage]
user=xstorage
group=xstorage

command=/home/xstorage/xstorage_env/bin/uwsgi --ini /home/xstorage/uwsgi.ini
environment=PATH="/home/xstorage/xstorage_env/bin:%(ENV_PATH)s"

directory=/home/xstorage
autostart=true
autorestart=true
startsecs=1
startretries=777
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
