# extra/supervisor.ini
[program:xstorage]
user=app_xstorage
group=app_xstorage

command=/home/app_xstorage/env/bin/uwsgi --ini /home/app_xstorage/uwsgi.ini
environment=PATH="/home/app_xstorage/env/bin:%(ENV_PATH)s"

directory=/home/xstorage
autostart=true
autorestart=true
startsecs=1
startretries=777
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5